# systemd/insta-timetable-daemon.service
[Unit]
Description=Insta Timetable Daemon (APScheduler, 07:00 Asia/Seoul)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/insta-timetable
EnvironmentFile=/opt/insta-timetable/.env
ExecStart=/opt/insta-timetable/scripts/run_daemon.sh
User=root
Group=root
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

