# systemd/insta-timetable-update.service
[Unit]
Description=Timetable update checker (NEIS hash compare)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/insta-timetable
EnvironmentFile=/opt/insta-timetable/.env
ExecStart=/opt/insta-timetable/venv/bin/python src/main_update.py
User=root
Group=root

[Install]
WantedBy=multi-user.target
